# nginix/supervisord.conf
[supervisord]
nodaemon=true

[program:python_app]
command=python /app/app.py
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true